# device-folio-configer/notify_system/systemd/notify_system.timer
[Unit]
Description=Notify System Timer
[Timer]
OnBootSec=30
OnUnitActiveSec=30
Unit=notify_system.service
[Install]
WantedBy=multi-user.target
